# deploy/roles/ovh-dyndns/files/dyndns_ovh.timer
[Unit]
Description=Update domain name adresses periodically

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
